# libs/nous_deploy/src/nous_deploy/vututi/config_templates/ututi_supervisord.conf
[program:{{ service.name }}]
directory={{ service.settings.instance_dir }}
command={{ service.settings.instance_dir }}/code/bin/paster serve {{ service.settings.instance_dir }}/release.ini
user={{ service.user }}
autostart=true
autorestart=true
redirect_stderr=true

stdout_logfile=/var/log/supervisor/{{service.name}}.log
stderr_logfile=/var/log/supervisor/{{service.name}}-error.log
